# .config/systemd/user/thunard.service
[Unit]
Description=Thunar file manager daemon

[Service]
ExecStart=/usr/sbin/thunar --daemon

[Install]
WantedBy=graphical.target

